# deploy/systemd/dream-motif-auto-sync.service
[Unit]
Description=Dream Motif Interpreter Google Docs Auto-Sync
After=network-online.target dream-motif-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/openclaw-you/workspace/Dream_Motif_Interpreter
EnvironmentFile=/srv/openclaw-you/workspace/Dream_Motif_Interpreter/.env
ExecStart=/srv/openclaw-you/workspace/Dream_Motif_Interpreter/.venv/bin/python -m app.auto_sync
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
